# deploy/systemd/cmem-server.service
[Unit]
Description=cmem-sync server
After=network.target
Documentation=https://github.com/bjarne/cmem-server

[Service]
Type=simple
User=cmem
Group=cmem
WorkingDirectory=/var/lib/cmem-server
ExecStart=/opt/cmem-server/cmem-server -c /etc/cmem-server.toml
Restart=on-failure
RestartSec=5

# 安全收紧
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/cmem-server
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictRealtime=true
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
